# otel/ghost/quadlets/otelcol.container
[Unit]
Description=OpenTelemetry Collector — aggregator (CNCF Incubating)
After=network-online.target loki.service prometheus.service

[Container]
Image=docker.io/otel/opentelemetry-collector-contrib:0.153.0
ContainerName=otelcol
Network=observability.network
Volume=%h/bluespeed/otel/config/otelcol-config.yaml:/etc/otelcol-contrib/config.yaml:ro,z
PublishPort=4317:4317
PublishPort=4318:4318
PublishPort=8888:8888

[Service]
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
