# otel/ghost/quadlets/perses.container
[Unit]
Description=Perses — dashboards (CNCF Sandbox)
After=network-online.target prometheus.service loki.service

[Container]
Image=docker.io/persesdev/perses:v0.53.1
ContainerName=perses
Network=observability.network
Volume=%h/bluespeed/otel/perses-data:/perses/db:z
Volume=%h/bluespeed/otel/config/perses-config.yaml:/etc/perses/config.yaml:ro,z
PublishPort=8082:8080

[Service]
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
